# rootfs/etc/s6-overlay/s6-rc.d/init-bind/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

mkdir -p /data/bind

if [[ ! -f /config/named.conf ]]; then
    cp /etc/bind/named.conf /config/named.conf
fi
