# rootfs/etc/s6-overlay/s6-rc.d/init-blocky/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

if [[ ! -f /config/blocky.yaml ]]; then
    cp /etc/blocky/config.yaml /config/blocky.yaml
fi
