# rootfs/etc/s6-overlay/s6-rc.d/init-dnsdist/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

if [[ ! -f /config/dnsdist.conf ]]; then
    cp /etc/dnsdist/dnsdist.conf /config/dnsdist.conf
fi
